# lobby-service/deploy/sts2-lobby.service.example
[Unit]
Description=STS2 Lobby Service
After=network.target
# Stop trying after 5 failures in 60s so deterministic boot errors (port
# already in use, missing files, bad permissions) surface as `failed` for
# an operator to investigate, instead of crash-looping every 3s and
# pummeling upstream peer discovery endpoints with repeat announces.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/sts2-lobby/lobby-service
# Parser defaults are explicit here; the later EnvironmentFile remains the
# operator override surface.
Environment="SERVER_CHAT_ENABLED=false"
Environment="SERVER_CHAT_RICH_CONTENT_ENABLED=true"
Environment="SERVER_CHAT_EMOJI_ENABLED=true"
Environment="SERVER_CHAT_ITEM_REFS_ENABLED=true"
Environment="ROOM_CHAT_V2_ENABLED=true"
Environment="ROOM_CHAT_COMBAT_REFS_ENABLED=true"
Environment="SERVER_CHAT_HISTORY_LIMIT=100"
Environment="SERVER_CHAT_HISTORY_TTL_HOURS=24"
Environment="SERVER_CHAT_SNAPSHOT_LIMIT=50"
Environment="SERVER_CHAT_MAX_PAYLOAD_BYTES=8192"
Environment="SERVER_CHAT_MAX_CONNECTIONS_PER_IP=10"
Environment="SERVER_CHAT_MAX_CONNECTIONS_TOTAL=500"
Environment="SERVER_CHAT_MAX_PENDING_TICKETS=2000"
Environment="SERVER_CHAT_TRUSTED_PROXY_CIDRS="
Environment="SERVER_CHAT_TICKET_REQUESTS_PER_MINUTE=20"
Environment="SERVER_CHAT_IP_MESSAGES_PER_MINUTE=30"
Environment="SERVER_CHAT_CONNECTION_BURST=5"
Environment="SERVER_CHAT_CONNECTION_REFILL_MS=2000"
Environment="SERVER_CHAT_SLOW_CLIENT_BYTES=262144"
Environment="MOD_SYNC_ENABLED=false"
Environment="MOD_SYNC_MAX_DESCRIPTORS=64"
Environment="MOD_SYNC_MAX_PAYLOAD_BYTES=65536"
EnvironmentFile=/opt/sts2-lobby/lobby-service/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=3
User=sts2
Group=sts2

[Install]
WantedBy=multi-user.target
