# lobby-service/deploy/sts2-lobby.service.example
[Unit]
Description=STS2 Lobby Service
After=network.target
# Stop trying after 5 failures in 60s so deterministic boot errors (port
# already in use, missing files, bad permissions) surface as `failed` for
# an operator to investigate, instead of crash-looping every 3s and
# pummeling upstream peer discovery endpoints with repeat announces.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/sts2-lobby/lobby-service
EnvironmentFile=/opt/sts2-lobby/lobby-service/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=3
User=sts2
Group=sts2

[Install]
WantedBy=multi-user.target
